# launchd/Library/LaunchAgents/com.jarrod.agentsview.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jarrod.agentsview</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/jarrod.folino/.local/bin/agentsview</string>
        <string>serve</string>
        <string>--no-browser</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/jarrod.folino/.local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>StandardOutPath</key>
    <string>/Users/jarrod.folino/.agentsview/launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jarrod.folino/.agentsview/launchd.err.log</string>
</dict>
</plist>
